# supervisord.conf
[program:main]
command = python maclient_web.py
directory = /home/bae/app
stdout_logfile = /home/bae/log/mac_web.log
stderr_logfile = /home/bae/log/mac_web_err.log

